# rc.local
su - ydeal -c "wetty -p 3000" >> /var/log/wetty.log 2>&1 &
